# poky/victor/meta-vicos/recipes-core/extra-conf/extra-conf/rsync/rsyncd.service
[Unit]
Description=rsync daemon
SourcePath=/usr/bin/rsync

[Service]
Restart=no
Type=forking
ExecStart=/usr/bin/rsync --daemon --config=/etc/rsyncd-victor.conf
