# poky/victor/meta-anki/recipes/anki-robot/files/anki-robot.target
[Unit]
Description=Anki Robotics System
Requires=multi-user-done.service
AllowIsolate=yes
